# opensolutions.service
[Unit]
Description=OpenSolutions Spring Boot Application
After=network.target mysql.service

[Service]
Type=simple
User=opensolutions
WorkingDirectory=/opt/opensolutions
ExecStart=/usr/bin/java -Xmx512m -jar /opt/opensolutions/opensolutions.jar
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
